# deploy/crontab
# m h  dom mon dow   command
02 7,8,9,10,12,14,16,19 * * * python /home/belem/webqa-dev/nightly_indicator/datahandler.py

